# system_build/game/rc.local
#!/bin/sh
su - game -c /usr/photon/bin/runmain &
